# linux/endeavouros/config/autostart/nvidia-perfmax.desktop
[Desktop Entry]
Type=Application
Name=Autoset Nvidia to Performance Mode
Exec="nvidia-settings" "-a" "[gpu:0]/GpuPowerMizerMode=1"
